# deploy/systemd/newscrawler-web.service
[Unit]
Description=Positive News Crawler web UI
After=network.target

[Service]
Type=simple
User=newscrawler
Group=newscrawler
WorkingDirectory=/opt/newscrawler
EnvironmentFile=/opt/newscrawler/.env
ExecStart=/opt/newscrawler/.venv/bin/python -m waitress --listen=127.0.0.1:8000 newscrawler.wsgi:application
Restart=always
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
