# deploy/systemd/newsagg-web.service
[Unit]
Description=Positive News Aggregator web UI
After=network.target

[Service]
Type=simple
User=newsagg
Group=newsagg
WorkingDirectory=/opt/newsagg
EnvironmentFile=/opt/newsagg/.env
ExecStart=/opt/newsagg/.venv/bin/python -m waitress --listen=127.0.0.1:8000 newsagg.wsgi:application
Restart=always
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

