# deploy/systemd/newscrawler-worker.service
[Unit]
Description=Positive News Crawler worker
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/newscrawler/newscrawler.env
ConditionPathExists=/var/lib/newscrawler/newscrawler.sqlite3
ConditionPathIsDirectory=/var/lib/newscrawler

[Service]
Type=simple
User=newscrawler
Group=newscrawler
WorkingDirectory=/opt/newscrawler
EnvironmentFile=/etc/newscrawler/newscrawler.env
Environment=PYTHONDONTWRITEBYTECODE=1
UMask=0007
ExecStartPre=/usr/bin/chmod 0660 /var/lib/newscrawler/newscrawler.sqlite3
ExecStart=/opt/newscrawler/.venv/bin/python manage.py runworker
Restart=always
RestartSec=10
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/newscrawler /var/log/newscrawler

[Install]
WantedBy=multi-user.target
